# tools/perkele/tests/syntax/list-reset.out.service
[Service]
Type=oneshot
ExecStart=
ExecStart=one
ExecStart=two
ExecStart=
ExecStart=three
ExecStart=four
