# tools/perkele/tests/syntax/list-reset.out.service
[Service]
Type=oneshot
ExecStart=
ExecStart=/bin/echo hello
ExecStart=/bin/echo bye
ExecStart=
ExecStart=/bin/echo HIYA
