# tests/tools/perkele/syntax/semicolons.out.service
[Service]
Type=oneshot
ExecStart=echo one ; echo "two two"
