# tools/perkele/tests/failure_handler/failure-handler@.out.service
[Unit]
Description=My failure handler for %i

[Service]
Type=oneshot
ExecStart=/usr/sbin/handler %i
