# tests/tools/perkele/simple/simple.service
[Unit]
Description=My Service

[Service]
ExecStart=/usr/bin/my-service

[Install]
WantedBy=multi-user.target
